# desktop/supervisor/conf.d/40_vncserver.conf
[program:vncserver]
command = bash /etc/sh/vncserver.sh
autostart = true
startsecs = 5
autorestart = true
startretries = 5
redirect_stderr = true
stdout_logfile_maxbytes = 20MB
stdout_logfile_backups = 5
stdout_logfile = /var/log/supervisor/vncserver.log
priority = 5
